# recipes-networking/chainball-ddns/files/duckdns.timer
[Unit]
Description=Update DuckDNS domain

[Timer]
OnUnitActiveSec=5m

[Install]
WantedBy=timers.target