# recipes-networking/scoreboard-https/files/dehydrated-renew.service
[Unit]
Description=Renew SSL certificate
After=Network.target

[Service]
Type=simple
ExecStart=__BINDIR__/dehydrated -c