# recipes-networking/scoreboard-https/files/dehydrated-renew.timer
[Unit]
Description=Run SSL renew monthly

[Timer]
OnCalendar=monthly
Persistent=true

[Install]
WantedBy=timers.target