# recipes-networking/scoreboard-https/files/dehydrated-register.service
[Unit]
Description=Register dehydrated on first boot
After=Network.target

[Service]
Type=simple
ExecStart=__BINDIR__/dehydrated-register.sh